# deployment/systemd/aas-api.service
[Unit]
Description=Autonomous Agent Stack API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/autonomous-agent-stack
EnvironmentFile=-/opt/autonomous-agent-stack/.env.linux
ExecStart=/opt/autonomous-agent-stack/.venv/bin/python -m uvicorn autoresearch.api.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
