# deployment/systemd/aas-housekeeper-mode-night.timer
[Unit]
Description=Run housekeeper night mode switch at 23:00

[Timer]
OnCalendar=*-*-* 23:00:00
Persistent=true

[Install]
WantedBy=timers.target
