# deployment/systemd/aas-housekeeper-mode-day.service
[Unit]
Description=Switch housekeeper to day_safe
After=aas-api.service

[Service]
Type=oneshot
WorkingDirectory=/opt/autonomous-agent-stack
EnvironmentFile=-/opt/autonomous-agent-stack/.env.linux
ExecStart=/opt/autonomous-agent-stack/scripts/housekeeper_api.sh mode-day
